# 7/etc/cron.daily/aide
# Run an AIDE system check every morning at 0500 local time.
0 5 * * * /usr/sbin/aide --check
